# sunrin-commute.conf
description "Gunicorn application server running sunrin-commute"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid root
setgid www-data

env PATH=/root/.python_virtual_envs/sunrin-commute/bin
chdir /root/vhosts/sunrin-commute
exec gunicorn --workers 3 --bind unix:sunrin-commute.sock -m 007 app:app
